# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[program:worker]
command=nameko run service --broker amqp://guest:guest@rabbit
autostart=true
autorestart=true
stopsignal=KILL
numprocs=1
stderr_logfile=/var/log/worker.err.log
stdout_logfile=/var/log/worker.out.log

[program:server]
command=flask run -h 0.0.0.0
autostart=true
autorestart=true
stopsignal=KILL
numprocs=1
stderr_logfile=/var/log/server.err.log
stdout_logfile=/var/log/server.out.log